# roles/wirerrd/templates/etc/systemd/system/wirerrd-export.timer.j2
[Unit]
Description=wirerrd export timer

[Timer]
Unit=wirerrd-export
OnBootSec=5m
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
